# open/software design/raspberry pi/systemctl/configure-gpio.service
[Unit]
Description=Configure GPIO
After=multi-user.target

[Service]
Type=simple
ExecStart=raspi-gpio set 14 a5

[Install]
WantedBy=multi-user.target